# rootfs/etc/services.d/fava/run
#!/usr/bin/execline -P
# load container environment.
with-contenv
s6-envdir /env/
exec s6-setuidgid abc /app/bin/fava